# skills/init/templates/claude-telegram-watchdog.timer.template
[Unit]
Description=Claude Telegram Listener Watchdog Timer ({{agent_name}})

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
