# skills/init/templates/claude-heartbeat.service.template
[Unit]
Description=Claude Code Heartbeat Runner ({{agent_name}})
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{working_directory}}
ExecStart=/bin/bash -c 'cat heartbeat-prompt.md | {{claude_path}} --model {{heartbeat_model}} -p --dangerously-skip-permissions'
TimeoutStartSec=3600

# Environment
Environment=HOME={{home_directory}}
Environment=PATH={{path}}
Environment=CLAUDE_CODE_EFFORT_LEVEL={{heartbeat_effort}}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-heartbeat
