# skills/init/templates/claude-telegram.service.template
[Unit]
Description=Claude Code Telegram Listener ({{agent_name}})
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{working_directory}}
ExecStart=/usr/bin/tmux new-session -ds {{tmux_session_name}} "{{claude_path}} --model {{listener_model}} --settings '{\"enabledPlugins\":{\"telegram@claude-plugins-official\":true}}' --channels plugin:telegram@claude-plugins-official --dangerously-skip-permissions --agent projects/remystack/agents/telegram-companion.md"
ExecStop=/usr/bin/tmux kill-session -t {{tmux_session_name}}
TimeoutStopSec=30

# Environment
Environment=HOME={{home_directory}}
Environment=PATH={{path}}
Environment=CLAUDE_CODE_EFFORT_LEVEL={{listener_effort}}

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-telegram

[Install]
WantedBy=default.target
